# docker_supervisor.conf
[supervisord]
nodaemon=true

[program:rails]
command=rails server -p 3001

[program:nginx]
command=nginx
